# bbc-on-ice.service
[Unit]
Description=bbc-on-ice

[Service]
ExecStart=/usr/bin/bbc-on-ice
