# bbc-on-ice.socket
[Socket]
ListenStream=[::1]:8000

[Install]
WantedBy=sockets.target
